# packaging/debian/postinst
#!/bin/sh
set -e
case "$1" in
  configure)
    # Create system user if missing.
    if ! getent passwd perseus-smarthome >/dev/null; then
      adduser --system --group --home /opt/raspberry-smarthome \
              --shell /usr/sbin/nologin --no-create-home \
              perseus-smarthome
    fi
    # Add to gpio for /dev/gpio* access.  usermod -aG is idempotent when
    # the user is already a member, so run it unconditionally.
    usermod -aG gpio perseus-smarthome
    chown -R perseus-smarthome:gpio /opt/raspberry-smarthome
    systemctl daemon-reload
    systemctl enable rpi-io-mcp.service
    systemctl start rpi-io-mcp.service
    ;;
esac
