# packaging/debian/prerm
#!/bin/sh
set -e
case "$1" in
  remove|deconfigure)
    if systemctl is-active --quiet rpi-io-mcp.service; then
      systemctl stop rpi-io-mcp.service
    fi
    if systemctl is-enabled --quiet rpi-io-mcp.service; then
      systemctl disable rpi-io-mcp.service
    fi
    ;;
  upgrade)
    # On upgrade: stop only; do not disable so the admin's choice is
    # preserved across the upgrade.
    if systemctl is-active --quiet rpi-io-mcp.service; then
      systemctl stop rpi-io-mcp.service
    fi
    ;;
esac
