# deploy/systemd/rpi-io-mcp.service
[Unit]
Description=Raspberry Pi I/O MCP Server
After=network.target

[Service]
Type=simple
User=pi
Group=gpio
WorkingDirectory=/opt/raspberry-smarthome
# Absolute path to uv: systemd's ExecStart binary lookup uses its own internal
# PATH (typically /usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin)
# and does NOT consult Environment="PATH=...", so a relative `uv run ...` here
# fails with status=203/EXEC when uv is installed under ~/.local/bin (the
# default for `pipx install uv`). The deploy script substitutes /home/pi/ with
# the actual deploy user's home before installing this unit.
Environment="PATH=/home/pi/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/pi/.local/bin/uv run rpi-io-mcp
Restart=on-failure
RestartSec=5
# Cap the SIGTERM-to-SIGKILL window so a hung asyncio shutdown cannot
# silently bypass GPIO teardown for the systemd default 90s.
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-io-mcp

# Allow the service to inherit GPIO device access from the gpio group.
# No special capabilities are required when the user is in the gpio group.

[Install]
WantedBy=multi-user.target
